# deploy/dmzcal-kiosk.service
[Unit]
Description=DMZ Dementia Calendar (Kiosk Mode)
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStartPre=/bin/sleep 3
ExecStart=/home/ubuntu/dmzcal/.venv/bin/dmzcal
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target default.target
